# roles/docker/templates/docker.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --graph="{{ docker_image }}"